# .config/containers/systemd/tor.container
[Unit]
Description=tor
Requires=servarr.pod
After=servarr.pod

[Container]
Image=docker.io/dockurr/tor:latest
ContainerName=tor
AutoUpdate=registry
Pod=servarr.pod

[Service]
Restart=always
TimeoutStartSec=60

